# docker/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/tmp/supervisord.pid
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:speedarr-backend]
command=uvicorn app.main:app --host 0.0.0.0 --port 9494
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DATABASE_URL="sqlite+aiosqlite:////data/speedarr.db"
